# AegisOSRepo/basic/services/aegis-license.service
[Unit]
Description=Aegis OS License Manager Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-license-daemon
Restart=on-failure
RestartSec=10
User=root
Group=root

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/aegis /etc/aegis

[Install]
WantedBy=multi-user.target
