# build-system/work/gamer/chroot/etc/systemd/system/aegis-render-orchestrator.service
[Unit]
Description=Aegis Render Orchestrator - Dual-GPU Split Rendering Daemon
Documentation=man:aegis-render-orchestrator(1)
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-render-orchestrator --daemon
Restart=on-failure
RestartSec=5
User=root
Group=root
StandardOutput=journal
StandardError=journal
Environment="XDG_RUNTIME_DIR=/run/user/0"

[Install]
WantedBy=graphical.target
