# aegis-os-freemium/overlay/etc/systemd/system/aegis-gaming-optimizer.service

[Unit]
Description=Aegis OS Gaming Optimizer - Freemium Edition
After=graphical-session.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/aegis-gaming-optimizer
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
