# aegis-os-freemium/overlay/etc/systemd/system/aegis-nouveau-optimizer.service

[Unit]
Description=Aegis OS Nouveau NVIDIA Driver Optimizer
After=graphical.target
Wants=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aegis-nouveau-optimizer
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
