# aegis-os-freemium/overlay/etc/systemd/system/aegis-gaming-optimizer.service

[Unit]
Description=Aegis OS Gaming Optimizer - System performance tuning for gaming
After=multi-user.target network.target
Wants=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/aegis-gaming-optimizer
RemainAfterExit=yes
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
