# aegis-os-freemium/overlay/etc/systemd/system/aegis-system-monitor.service

[Unit]
Description=Aegis OS System Monitor - Freemium Edition
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/aegis-system-monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
