# aegis-os-freemium/overlay/etc/systemd/system/aegis-kernel-interface.service

[Unit]
Description=Aegis OS Kernel Interface - Kernel module management and integration
After=multi-user.target network.target systemd-modules-load.service
Wants=systemd-modules-load.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/aegis-kernel-interface --load
ExecStart=/usr/local/bin/aegis-kernel-interface --status
RemainAfterExit=yes
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
