# build-system/work/basic/chroot/etc/systemd/system/aegis-security-suite.service
[Unit]
Description=Aegis Security Suite - Integrated security monitoring for Basic edition
Documentation=https://aegis-os.com/docs/security-suite
After=multi-user.target network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/aegis-security-suite --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10

TimeoutStartSec=30
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=aegis-security-suite

PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=false
ReadWritePaths=/var/lib/aegis /var/log/aegis /etc/aegis

Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"

[Install]
WantedBy=multi-user.target
