# build-system/archiso/profiles/releng/airootfs/etc/systemd/system/aegis-live.service
[Unit]
Description=Aegis OS Live Environment Setup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/aegis-live-setup

[Install]
WantedBy=multi-user.target
