# AegisOSRepo/gamer/services/aegis-streamforge.service
[Unit]
Description=Aegis OS StreamForge Capture
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-streamforge-daemon
Restart=on-failure
RestartSec=5
User=aegis

[Install]
WantedBy=graphical.target
