# AegisOSRepo/gamer/services/aegis-display-optimizer.service
[Unit]
Description=Aegis OS Display Optimization Service
After=multi-user.target graphical.target
Wants=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-display-optimizer-daemon
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=graphical.target
