# AegisOSRepo/gamer/services/aegis-latency-fastpath.service
[Unit]
Description=Aegis OS Latency FastPath
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/aegis-latency-fastpath enable
ExecStop=/usr/local/bin/aegis-latency-fastpath disable
User=root

[Install]
WantedBy=multi-user.target
