# AegisOSRepo/gamer/services/aegis-vram-balancer.service
[Unit]
Description=Aegis OS VRAM Balancer
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-vram-balancer-daemon
Restart=on-failure
User=root

[Install]
WantedBy=graphical.target
