# AegisOSRepo/gamer/services/aegis-shader-precache.service
[Unit]
Description=Aegis OS Shader Pre-Cache Engine
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-shader-precache-daemon
Restart=on-failure
User=aegis
Nice=19
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
