# AegisOSRepo/gamer/services/aegis-thermal-guard.service
[Unit]
Description=Aegis OS Thermal Guard
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-thermal-guard-daemon
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
