# build-system/overlays/gamer/etc/udev/rules.d/99-aegis-external-games.rules
# Aegis OS External Game Media udev Rules
# Automatically mount USB drives and SD cards for game launching

# USB Storage Devices - trigger game media detection
ACTION=="add", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", ENV{DEVTYPE}=="partition", \
    RUN+="/usr/local/bin/aegis-external-media-handler mount %k"

ACTION=="remove", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", ENV{DEVTYPE}=="partition", \
    RUN+="/usr/local/bin/aegis-external-media-handler unmount %k"

# SD Card Devices (mmc)
ACTION=="add", SUBSYSTEM=="block", KERNEL=="mmcblk[0-9]p[0-9]*", \
    RUN+="/usr/local/bin/aegis-external-media-handler mount %k"

ACTION=="remove", SUBSYSTEM=="block", KERNEL=="mmcblk[0-9]p[0-9]*", \
    RUN+="/usr/local/bin/aegis-external-media-handler unmount %k"

# NVMe external drives (USB enclosures)
ACTION=="add", SUBSYSTEM=="block", KERNEL=="nvme[0-9]n[0-9]p[0-9]*", ENV{ID_BUS}=="usb", \
    RUN+="/usr/local/bin/aegis-external-media-handler mount %k"

ACTION=="remove", SUBSYSTEM=="block", KERNEL=="nvme[0-9]n[0-9]p[0-9]*", ENV{ID_BUS}=="usb", \
    RUN+="/usr/local/bin/aegis-external-media-handler unmount %k"
