# build-system/overlays/gamer/etc/systemd/user/aegis-wallpaper.service
[Unit]
Description=Aegis Wallpaper Engine - Dynamic Wallpapers
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aegis-wallpaper daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
