# build-system/overlays/gamer/etc/systemd/user/aegis-gaming-mode.target
[Unit]
Description=Aegis Gaming Mode - All Gaming Optimizations Active
Wants=aegis-ram-guardian.service
Wants=aegis-latency-fastpath.service
Wants=aegis-netboost.service
Wants=aegis-shader-precache.service
Wants=aegis-audio-lowlatency.service
Wants=aegis-thermal-guard.service
Wants=aegis-vram-balancer.service

[Install]
WantedBy=graphical-session.target
